# recipes-rtd-mcu/rtd-mcu-bin/files/mcu_sbl_mode.service
[Unit]
Description=To Reset the MCU and to run the user application in MCU

[Service]
Type=oneshot
ExecStart=/usr/bin/mcu_reset.sh

[Install]
WantedBy=multi-user.target
